# .config/production/supervisord.conf
[program:nginx]
command=nginx

[program:uwsgi]
command=/root/.pyenv/versions/ticket/bin/uwsgi --ini /srv/ticket/.config/production/uwsgi.ini

