# .config/production/supervisord.conf
[program:nginx]
command=nginx

[program:uwsgi]
command=uwsgi --ini /srv/ticket/.config/production/uwsgi.ini

[program:cron]
command = cron -f -L 15
autostart = true
autorestart = true